# on branch release.ini
[program:kargopay]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/port_simplified/artisan queue:work --queue=default --sleep=3  --timeout=40000 --max-time=40000 --max-jobs=1
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=16
redirect_stderr=true
stdout_logfile=/var/www/html/port_simplified/storage/logs/worker.log

[program:kpa_payment_confirmation_responses]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/port_simplified/artisan rabbitmq:consume kpa_payment_confirmation_responses
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/port_simplified/storage/logs/worker.log
stopwaitsecs=3600

[program:kpa_payment_ref_responses]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/port_simplified/artisan rabbitmq:consume kpa_payment_ref_responses
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/port_simplified/storage/logs/worker.log
stopwaitsecs=3600

[program:kpa_invoices_responses]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/port_simplified/artisan rabbitmq:consume kpa_invoices_responses
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/port_simplified/storage/logs/worker.log
stopwaitsecs=3600

[program:kpa_query_invoice_responses]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/port_simplified/artisan rabbitmq:consume kpa_query_invoice_responses
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/html/port_simplified/storage/logs/worker.log
stopwaitsecs=3600




* * * * * cd /var/www/html/port_simplified && php artisan schedule:run >> /dev/null 2>&1

SELECT * FROM public.users ORDER BY id DESC LIMIT 2;